# deployed-configs/staging/systemd/logans3d-staging.service
[Unit]
Description=Logan's 3D Creations Go Application (Staging)
After=network.target

[Service]
Type=simple
User=apprunner
Group=apprunner
EnvironmentFile=/etc/logans3d-staging/environment
WorkingDirectory=/home/apprunner/sites/logans3d-staging
ExecStart=/home/apprunner/sites/logans3d-staging/logans3d
Restart=always
RestartSec=10
StandardOutput=append:/var/log/logans3d-staging/logans3d-staging.log
StandardError=append:/var/log/logans3d-staging/logans3d-staging.log

[Install]
WantedBy=multi-user.target